# ops/systemd/finance-pi-backup.timer
[Unit]
Description=Run verified finance-pi backup weekly

[Timer]
OnCalendar=Sun *-*-* 03:30:00 Asia/Seoul
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
